# bin404-cleanup/bin404-cleanup.timer
[Unit]
Description=Run bin404-cleanup every hour

[Timer]
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
